# rack0_master/maas/rc.local
mount --bind /usr/src/maas_hooks/hooks.py /snap/maas/current/lib/python3.6/site-packages/metadataserver/builtin_scripts/hooks.py
